# src/pa/packaging/systemd/pa-server.service.template
[Unit]
Description=PA Server ({{PA_INSTANCE_NAME}})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{PA_BIN}} serve
Restart=on-failure
RestartSec=5
TimeoutStopSec=300
LimitNOFILE=8192
MemoryMax=4G

{{ENV_LINES}}

StandardOutput=append:{{PA_LOG_DIR}}/server.log
StandardError=append:{{PA_LOG_DIR}}/server.err.log

[Install]
WantedBy=default.target
